# ptt-bridge/install/macos/com.council.ptt-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.council.ptt-bridge</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/usr/local/lib/council-ptt-bridge/dist/ptt-bridge/src/index.js</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PTT_BRIDGE_HOST</key>
    <string>127.0.0.1</string>
    <key>PTT_BRIDGE_PORT</key>
    <string>8765</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/var/log/council-ptt-bridge.log</string>
  <key>StandardErrorPath</key>
  <string>/var/log/council-ptt-bridge.err.log</string>
</dict>
</plist>
